# sound/snapcast/files/snapserver.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2024 OpenWrt.org

START=99
BIN=/usr/bin/snapserver
PID=/var/run/snapserver.pid
SSD=start-stop-daemon

start() {
        ${SSD} -p ${PID} -S -x ${BIN} -- -P ${PID}
}

stop() {
        ${SSD} -p ${PID} -K -s SIGINT
}
